# com.ai-relay.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-relay.tunnel</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/ssh</string>
        <string>-i</string>
        <string>/Users/manu/.ssh/contabo_key</string>
        <string>-R</string>
        <string>4003:localhost:4003</string>
        <string>-R</string>
        <string>5050:localhost:5050</string>
        <string>-N</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=no</string>
        <string>admin@89.117.61.193</string>
    </array>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/manu/Documents/DEV/ai-relay-server/tunnel.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/manu/Documents/DEV/ai-relay-server/tunnel.error.log</string>
</dict>
</plist>